# sys-analytics/.config/systemd/user/sys-analytics.timer
[Unit]
Description=System analytics collector timer

[Timer]
OnBootSec=10s
OnUnitActiveSec=30s
AccuracySec=5s

[Install]
WantedBy=timers.target
