# sys-analytics/.config/systemd/user/sys-analytics-dash.service
[Unit]
Description=System analytics dashboard generator

[Service]
Type=oneshot
ExecStart=%h/dotfiles/sys-analytics/dashboard.sh
Nice=19
